# deploy/ec2/ce/web.service
[Unit]
Description=Web Server
After=network.target

[Service]
Type=simple
User=ubuntu

WorkingDirectory=/home/ubuntu/app
Environment="NODE_ENV=production"
EnvironmentFile=/home/ubuntu/app/.env
RestartSec=1
ExecStart=/usr/bin/pnpm --prefix /home/ubuntu/app start:prod
Restart=always

[Install]
WantedBy=multi-user.target